# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=5
loglevel=info
pidfile=/var/run/supervisord.pid
user=root

; ─── Nginx reverse proxy ──────────────────────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/nginx.stdout.log
stderr_logfile=/var/log/supervisor/nginx.stderr.log

; ─── Primary Backend (Auth + Dashboard API) — port 3000 ─────────────────────
[program:primary-backend]
command=bun run apps/primary-backend/src/index.ts
directory=/app
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/primary-backend.stdout.log
stderr_logfile=/var/log/supervisor/primary-backend.stderr.log
environment=PORT="3000"

; ─── API Gateway (LLM proxy) — port 4000 ────────────────────────────────────
[program:api-gateway]
command=bun run apps/api-backend/src/index.ts
directory=/app
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/api-gateway.stdout.log
stderr_logfile=/var/log/supervisor/api-gateway.stderr.log
environment=PORT="4000"
